# scripts/workhorse-station.service
[Unit]
Description=Workhorse Station - AI Dev Management Web App
After=network.target

[Service]
Type=simple
WorkingDirectory=PROJECT_ROOT
EnvironmentFile=-PROJECT_ROOT/.env
ExecStart=/usr/bin/env node PROJECT_ROOT/apps/api/dist/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
